# app/supervisor/chronograf.conf
[program:chronograf]
user = chronograf
command = /usr/bin/chronograf -b /data/chronograf/chronograf-v1.db
autostart = true
autorestart = true
startsecs = 5
