# app/supervisor/chronograf.conf
[program:chronograf]
user = chronograf
command = /usr/bin/chronograf -b /data/chronograf/chronograf-v1.db --port=8080
autostart = true
autorestart = true
startsecs = 5
