# app/supervisor/kapacitor.conf
[program:kapacitor]
user = kapacitor
command = /usr/bin/kapacitord -config /etc/kapacitor/kapacitor.conf
autostart = true
autorestart = true
startsecs = 5